# deploy/anilistgen.timer
[Unit]
Description=Run anilistgen daily for seasonal anime MDBList sync
Documentation=https://github.com/calmcacil/anilistgen

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
