# deploy/animelistgen.service
[Unit]
Description=animelistgen — seasonal anime MDBList sync
Documentation=https://github.com/calmcacil/animelistgen
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/animelistgen -c /etc/animelistgen/animelistgen.yaml
User=animelistgen
Group=animelistgen
Restart=no
